# com.user.batch-timer.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>KeepAlive</key>
    <false/>
    <key>Label</key>
    <string>com.user.batch-timer.daily</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>REMINDER_TITLE</key>
      <string>Batch Timer</string>
      <key>REMINDER_MESSAGE</key>
      <string>20:00 のリマインダです</string>
      <key>REMINDER_SOUND</key>
      <string>default</string>

      <key>BATCH_TIMER_STAMP_DIR</key>
      <string>/Users/rchaser53/batch-timer/logs/stamps</string>
      <key>BATCH_TIMER_ONCE_LOG_FILE</key>
      <string>/Users/rchaser53/batch-timer/logs/once-per-day.log</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/rchaser53/batch-timer/scripts/reminder-check.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/rchaser53/batch-timer/logs/stderr.log</string>
    <key>StandardOutPath</key>
    <string>/Users/rchaser53/batch-timer/logs/stdout.log</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>20</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StartInterval</key>
    <integer>60</integer>
  </dict>
</plist>