# fetch-usd-jpy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>EnvironmentVariables</key>
    <dict>
      <key>BATCH_TIMER_CATCHUP</key>
      <string>1</string>
      <key>REMINDER_TITLE</key>
      <string>Batch Timer</string>
      <key>REMINDER_MESSAGE</key>
      <string/>
      <key>REMINDER_SOUND</key>
      <string>default</string>
      <key>REMINDER_MODE</key>
      <string>alert</string>
    </dict>
    <key>KeepAlive</key>
    <false/>
    <key>Label</key>
    <string>com.user.batch-timer.sample.every-5min</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/rchaser53/batch-timer/scripts/sample-task.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/rchaser53/batch-timer/logs/sample.stderr.log</string>
    <key>StandardOutPath</key>
    <string>/Users/rchaser53/batch-timer/logs/sample.stdout.log</string>
    <key>StartInterval</key>
    <integer>307</integer>
  </dict>
</plist>