# web-test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>EnvironmentVariables</key>
    <dict>
      <key>REMINDER_TITLE</key>
      <string>Batch Timer</string>
      <key>REMINDER_MESSAGE</key>
      <string/>
      <key>REMINDER_SOUND</key>
      <string>default</string>
      <key>REMINDER_MODE</key>
      <string>alert</string>
    </dict>
    <key>KeepAlive</key>
    <false/>
    <key>Label</key>
    <string>com.user.batch-timer.web-test.every-10min</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/rchaser53/batch-timer/scripts/web-test-task.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/rchaser53/batch-timer/logs/web-test.stderr.log</string>
    <key>StandardOutPath</key>
    <string>/Users/rchaser53/batch-timer/logs/web-test.stdout.log</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>20</integer>
      <key>Minute</key>
      <integer>26</integer>
    </dict>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>WorkingDirectory</key>
    <string>/Users/rchaser53/batch-timer</string>
  </dict>
</plist>